# ansible/roles/service/files/service/rateupdate.service
[Unit]
Description=Update AMD/USD rate 

[Service]
ExecStart=/bin/bash /opt/rateupdate.sh

[Install]
WantedBy=multi-user.target